# infra/systemd/orchestrator-backup.service
[Unit]
Description=Orchestrator PostgreSQL daily backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/codegen_orchestrator/infra/scripts/backup-db.sh
EnvironmentFile=/opt/codegen_orchestrator/.env
User=deploy
Group=deploy

[Install]
WantedBy=multi-user.target
