# infra/systemd/orchestrator-backup.timer
[Unit]
Description=Run orchestrator DB backup daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
